# conf/local/etc/supervisord.conf
;/etc/supervisord.conf
[unix_http_server]
file = /var/run/supervisord.sock
chmod = 0777
chown= ttwait

[inet_http_server]
# Web管理界面设定
port=127.0.0.1:9001
username = 
password = 

[supervisorctl]
; 必须和'unix_http_server'里面的设定匹配
serverurl = unix:///var/run/supervisord.sock

[supervisord]
logfile=/var/log/supervisord/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root                 ; (default is current user, required if root)
childlogdir=/var/log/supervisord/            ; ('AUTO' child log dir, default $TEMP)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; 管理的单个进程的配置，可以添加多个program




[program:nginx]
command=/usr/local/sbin/nginx 
autostart = true
user = root
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisord/nginx.log

[program:mysqld]
command=/usr/local/mysql/bin/mysqld
autostart = true
user = ttwait
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisord/mysqld.log

[program:redis]
command=/usr/local/bin/redis-server
autostart = true
user = ttwait
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisord/redis-server.log

[program:mongod]
command=/usr/local/bin/mongod --rest
autostart = true
user = ttwait
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisord/mongod.log


[program:diycelerd]
command=/Users/ttwait/work/code/bitbucket/diy/env/bin/python run.py worker --loglevel=INFO --concurrency=1 -n speedycustom.com --autoreload -E --config=celeryconfig --maxtasksperchild=10000 --time-limit=3600 --soft-time-limit=3600
autostart = true
directory = /Users/ttwait/work/code/bitbucket/diy/task_queue
user = ttwait
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisord/diycelerd.log


[program:diyflower]
command=/Users/ttwait/work/code/bitbucket/diy/env/bin/python run.py flower --address=127.0.0.1 --port=5555 --broker=redis://localhost:6379/0
autostart = true
directory = /Users/ttwait/work/code/bitbucket/diy/task_queue
user = ttwait
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisord/diyflower.log


[program:searchd]
command=/usr/local/coreseek/bin/searchd -c conf/zen_dev/sphinx/sku.conf --nodetach
autostart = false
startsecs = 3
user = webuser
environment=DYLD_LIBRARY_PATH=/usr/local/mysql/lib
directory = /Users/ttwait/work/code/mfhui/
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/supervisord/searchd.log


